# launchd/com.zachkysar.jpr-pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.zachkysar.jpr-pipeline</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/zachkysar/projects/utils/jpr-pipeline/jpr_pipeline.py</string>
  </array>

  <key>WatchPaths</key>
  <array>
    <string>/Users/zachkysar/Library/Mobile Documents/iCloud~com~openplanetsoftware~just-press-record/Documents</string>
  </array>

  <key>StartInterval</key>
  <integer>60</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/jpr-pipeline.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/jpr-pipeline-error.log</string>
</dict>
</plist>
